# asset/systemd/buh-node-sweep.timer
[Unit]
Description=Periodic buh TTL sweep of expired envelopes

[Timer]
# Envelopes also expire lazily on read; this sweep reclaims space for queues that are never
# pulled again. Hourly is ample for a store-and-forward relay.
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m
Unit=buh-node-sweep.service

[Install]
WantedBy=timers.target
