# asset/systemd/buh-node-sweep.service
[Unit]
Description=buh TTL sweep — delete expired envelopes from the relay datastore

[Service]
Type=oneshot
User=buh
Group=buh
ExecStart=/usr/local/bin/buh-cli --db-path /var/lib/buh/relay.db sweep

# Same hardening posture as the node service; only the datastore is writable.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/buh
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
